# docker/gunicorn.conf
[program:gunicorn]
command=newrelic-admin run-program gunicorn --config=/code/gunicorn.py steamworkshopdownloader:app
directory=/code
